# units/lazy_mount.service
[Unit]
Description=Mount /opt/usr

[Service]
Type=oneshot
SmackProcessLabel=System::Privileged
ExecStart=/usr/bin/mount-user.sh
SuccessExitStatus=0
